# systemd/update-goes-fd-web.service
[Unit]
Description=Publish latest GOES-19 Full Disk images for web UI

[Service]
Type=oneshot
User=root
Group=root
ExecStart=/usr/local/bin/update_goes_fd_web.sh
